# phoneresponse.service
[Unit]
Description=PhoneResponder

[Service]
ExecStart=/home/jpw/phoneresponse
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target